# scripts/drive-cleanup-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bmc.drive-cleanup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/matias/Panelin calc loca/Calculadora-BMC/scripts/drive-cleanup-automated.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/matias/.cache/drive-cleanup/launchd-out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/matias/.cache/drive-cleanup/launchd-err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>THRESHOLD_GB</key>
    <string>8</string>
    <key>LOG_DIR</key>
    <string>/Users/matias/.cache/drive-cleanup</string>
  </dict>
</dict>
</plist>
